# src/files/debian/etc/profile.d/hugo_shell.sh
#!/bin/sh

# Change layout of shell prompt
PS1="\e[1;32mhugo\e[m:\e[1;34m\w\e[m\$ "

# Create alias for `hugo`
alias h="hugo"
