# src/files/debian/etc/profile.d/hugo.sh
#!/bin/sh

# Environment configuration triggered by both shell and hugo script.
